# config/etc/service/f5-icontrol-trusted-devices-express/run
#!/bin/sh

cd /var/lib/f5-icontrol-trusted-devices
# Start express service
export PORT=8105
exec /usr/bin/node index.js 