# config/etc/service/apache/finish
#!/bin/sh

NAME=httpd
PID=`pgrep $NAME`
echo $PID > log
date >> log
wait $PID