# config/etc/service/restjavad/finish
#!/bin/sh

NAME=java
PID=`pgrep $NAME`
echo $PID > log
date >> log
wait $PID
