# config/etc/service/nginx/finish
#!/bin/sh

PID=`pgrep "nginx: master"`
echo $PID > log
date >> log
wait $PID
