# config/etc/service/unit/finish
#!/bin/sh

PID=`pgrep /usr/sbin/unitd`
echo $PID > log
date >> log
wait $PID
